# debian/cblues-xfce-livesettings.service
[Unit]
Description=Xfce cheatcodes in live media
ConditionPathExists=/usr/share/cblues-xfce-settings/
After=cblues-common-setting.service

[Service]
Type=oneshot
ExecStart=/usr/share/cblues-xfce-settings/xfce-cheatcodes start

[Install]
WantedBy=graphical.target
