# debian/farewell-xfce-livesettings.service
[Unit]
Description=Xfce cheatcodes in live media
ConditionPathExists=/usr/share/farewell-xfce-settings/
After=farewell-common-setting.service

[Service]
Type=oneshot
ExecStart=/usr/share/farewell-xfce-settings/xfce-cheatcodes start

[Install]
WantedBy=graphical.target
