# debian/wintersky-xfce-livesettings.service
[Unit]
Description=Xfce cheatcodes in live media
ConditionPathExists=/usr/share/wintersky-xfce-settings/
After=wintersky-common-setting.service

[Service]
Type=oneshot
ExecStart=/usr/share/wintersky-xfce-settings/xfce-cheatcodes start

[Install]
WantedBy=graphical.target
